# dashboard/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
command=gunicorn --bind 127.0.0.1:5000 --timeout 600 api_server:app
directory=/home/site/wwwroot
autostart=true
autorestart=true
stdout_logfile=/tmp/api.log
stderr_logfile=/tmp/api_error.log

[program:dashboard]
command=streamlit run Home.py --server.port=8501 --server.address=127.0.0.1 --server.headless=true
directory=/home/site/wwwroot
autostart=true
autorestart=true
stdout_logfile=/tmp/dashboard.log
stderr_logfile=/tmp/dashboard_error.log

[program:nginx]
command=nginx -c /home/site/wwwroot/nginx.conf -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/tmp/nginx.log
stderr_logfile=/tmp/nginx_error.log
